# services/dns/run
#!/usr/bin/env bash

set -Eeuo pipefail

for domain in ${MAISON_DOMAINS}; do
    echo "$MAISON_LOCAL_IP $domain" >> /etc/hosts
done

dnsmasq -d
