# services/dns/run
#!/usr/bin/env bash

set -Eeuo pipefail

for service in ${MAISON_WEB_SERVICES}; do
    echo "$MAISON_LOCAL_IP $service.$MAISON_WEB_DOMAIN" >> /etc/hosts
done

dnsmasq -d
