# services/vcs/run
#!/usr/bin/env bash

set -Eeuo pipefail

envsubst < /maison/app.ini.envtpl > /data/gitea/conf/app.ini

chown "$(whoami):$(groups | awk '{ print $1 }')" -R /data

/usr/bin/entrypoint "$@"
