# services/hab/run
#!/usr/bin/env bash

set -Eeuo pipefail

templates="/openhab/conf/services/jdbc.cfg"
for template in $templates; do
    envsubst < $template.envtpl > $template
done

/entrypoint.sh "$@"
