# services/consul/run
#!/usr/bin/env bash

set -Eeuo pipefail

esh -o "/consul/config/consul.json" "/consul/config/consul.json.esh"

consul agent -config-dir=/consul/config
