# services/web/run
#!/usr/bin/env bash

set -e

export MAISON_HOST_IP=$(ip route show | awk '/default/ {print $3}')

templates=/etc/nginx/conf.d/*.conf*.esh
for template in $templates; do
    esh -o "${template::-4}" "$template"
done

# Obtain HTTPS certificates.
nginx
certbot_args="--nginx -n --agree-tos --email bart@mynameisbart.com"
for domain in ${MAISON_DOMAINS}
do
    certbot_args="$certbot_args -d $domain"
done
# If HTTPS is not enabled explicitly, use dry runs and allow failures, as those may be caused by non-public domains
# (*.local), which Let's Encrypt rejects.
if [ ${MAISON_WEB_PROTOCOL} == 'https' ]; then
    certbot $certbot_args
else
    certbot certonly --dry-run $certbot_args || true
fi;
nginx -s quit

# Re-launch nginx in the foreground.
nginx -g 'daemon off;'
