# services/mail/conf/supervisor/postfix.ini
[program:postfix]
command=postfix start
