# services/mailer/conf/runit/postfix/run
#!/usr/bin/env sh

esh -o /etc/postfix/main.cf /etc/postfix/main.cf.esh

postfix start
